# roles/singbox_tproxy/templates/etc/systemd/system/sing-box-reload.path.j2
[Unit]
Description=Watch {{ singbox_etc_dir }}/config.json for changes

[Path]
PathChanged={{ singbox_etc_dir }}/config.json
Unit=sing-box-reload.service

[Install]
WantedBy=multi-user.target
