# roles/singbox_config/templates/etc/systemd/system/sing-box-config-updater.service.j2
[Unit]
Description=sing-box-config updater

[Service]
Type=oneshot
WorkingDirectory={{ singbox_state_dir }}
ExecStart={{ singbox_venv_dir }}/bin/sing-box-config --output {{ singbox_etc_dir }}/config.json
